# deploy/systemd/ha-predictor-training.service
[Unit]
Description=HA Intent Predictor - Model Training Service
Documentation=https://github.com/yourusername/ha-intent-predictor
After=network.target docker.service ha-intent-predictor.service
Wants=network.target
Requires=docker.service
PartOf=ha-intent-predictor.service

[Service]
Type=exec
User=root
Group=root
WorkingDirectory=/opt/ha-intent-predictor
Environment=PATH=/opt/ha-intent-predictor/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONPATH=/opt/ha-intent-predictor
Environment=LOG_LEVEL=INFO
Environment=SERVICE_NAME=training
EnvironmentFile=-/opt/ha-intent-predictor/.env

# Model training service startup
ExecStart=/opt/ha-intent-predictor/venv/bin/python -m src.learning.main
ExecReload=/bin/kill -HUP $MAINPID

# Process management
Restart=always
RestartSec=15
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=60

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/ha-intent-predictor /var/log/ha-intent-predictor /tmp

# Resource limits (most intensive service)
LimitNOFILE=65536
LimitNPROC=32768
MemoryAccounting=yes
MemoryHigh=4G
MemoryMax=5G
CPUAccounting=yes
CPUQuota=300%

# Standard output to journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ha-predictor-training

[Install]
WantedBy=multi-user.target